# builder/files/etc/rc.local
#!/bin/bash


if [ -f /aafirstboot ]; then /aafirstboot start ; fi

exit 0
