# builder/files/etc/systemd/system/docker.service.d/10-execstart.conf
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon -H fd:// --log-driver=journald --storage-driver=overlay --icc=false --iptables=true
